# service/node/run
#!/bin/sh
echo starting
exec static /home/cloud/git/blackjack/static/
